# crontab/root
*/1 * * * * cd /code && python3 /code/app.py
